# packages/genocide/genocide-38.tar.gz/genocide-38/files/genocide
#!/bin/ksh
# This file is placed in the Public Domain.

daemon="/usr/local/bin/genocide"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
            groupadd genocide
            useradd -b /var/lib -d /var/lib/genocide -m -g genocide genocide 2>&1 &
}

rc_stop() {
            pkill -f genocide
}

rc_cmd $1
