# packages/mopidris/mopidris-0.9.0.tar.gz/mopidris-0.9.0/mopidris.service
[Unit]
Description=Mopidy MPRIS2 D-Bus service
Documentation=https://codeberg.org/sph/mopidris
Wants=dbus.socket

[Service]
ExecStart=python3 mopidris
Restart=on-success

[Install]
WantedBy=multi-user.target